# oracle/crediclip-worker.service
[Unit]
Description=CrediClip Queue Worker
After=network-online.target ollama.service
Wants=network-online.target
Requires=ollama.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/CrediClip
EnvironmentFile=-/etc/default/crediclip-worker
Environment=PYTHONPATH=/home/ubuntu/CrediClip
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc 'source /home/ubuntu/CrediClip/.venv/bin/activate && python scripts/queue_worker.py --api-base "${API_BASE:-https://crediclip-axraza-msba.fly.dev}" --worker-id "${WORKER_ID:-oracle-worker-1}" --ingest-mode "${INGEST_MODE:-rich}" --analysis-mode "${ANALYSIS_MODE:-local}" --poll-sec "${POLL_SEC:-3}" --timeout-sec "${TIMEOUT_SEC:-180}"'
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
